# deduper.service
[Unit]
Description=DeDuper Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/your/script.py
Restart=on-failure

[Install]
WantedBy=multi-user.target